# silabs_examples/unify-matter-bridge/linux/debconf/prerm
#!/bin/sh

if type "systemctl" > /dev/null; then
    systemctl stop uic-mb.service
    systemctl disable uic-mb.service
    rm -rf /var/lib/uic-mb/*
fi